# deploy/sub2api.service
[Unit]
Description=Sub2API Admin Plus - operations automation extension
Documentation=https://github.com/openrelayllm/sub2api-admin-plus
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=simple
User=sub2api
Group=sub2api
WorkingDirectory=/opt/sub2api-admin-plus
ExecStart=/opt/sub2api-admin-plus/sub2api-admin-plus
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sub2api-admin-plus

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/opt/sub2api-admin-plus
ReadWritePaths=/etc/sub2api-admin-plus

# Environment - Server configuration
# Modify these values to change listen address and port
Environment=GIN_MODE=release
Environment=DATA_DIR=/etc/sub2api-admin-plus
Environment=SERVER_HOST=0.0.0.0
Environment=SERVER_PORT=8080
Environment=ADMIN_PLUS_PROXY_MIHOMO_BINARY_PATH=/opt/sub2api-admin-plus/bin/mihomo

[Install]
WantedBy=multi-user.target
